# contrib/init/bitcoind.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SOLARID_BIN="/usr/bin/solarid"
env SOLARID_USER="solari"
env SOLARID_GROUP="solari"
env SOLARID_PIDDIR="/var/run/solarid"
# upstart can't handle variables constructed with other variables
env SOLARID_PIDFILE="/var/run/solarid/solarid.pid"
env SOLARID_CONFIGFILE="/etc/solari/solari.conf"
env SOLARID_DATADIR="/var/lib/solarid"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # solarid will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SOLARID_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run solarid."
        echo "The setting must appear in $SOLARID_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SOLARID_PIDDIR"
    chmod 0755 "$SOLARID_PIDDIR"
    chown $SOLARID_USER:$SOLARID_GROUP "$SOLARID_PIDDIR"
    chown $SOLARID_USER:$SOLARID_GROUP "$SOLARID_CONFIGFILE"
    chmod 0660 "$SOLARID_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SOLARID_PIDFILE" \
    --chuid $SOLARID_USER:$SOLARID_GROUP \
    --exec "$SOLARID_BIN" \
    -- \
    -pid="$SOLARID_PIDFILE" \
    -conf="$SOLARID_CONFIGFILE" \
    -datadir="$SOLARID_DATADIR" \
    -disablewallet \
    -daemon

